# paas-ce/paas/examples/supervisord.d/appengine.ini
[program:appengine]
command=/opt/dev-paas/paas-runtime/appengine/bin/gunicorn wsgi --bind 0.0.0.0:8000 -k gevent -w 4 -n appengine --access-logfile - --error-logfile -
directory=/opt/dev-paas/paas-ce/paas/appengine
environment=PYTHONPATH=/opt/dev-paas/paas-runtime/appengine,BK_ENV="production",PAAS_LOGGING_DIR="/opt/dev-paas/paas-runtime/logs"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/dev-paas/paas-runtime/logs/paas_appengine.log
