# paas-ce/paas/examples/supervisord.d/esb.ini
[program:esb]
command=/opt/opsany/.runtime/esb/bin/gunicorn wsgi --bind 0.0.0.0:8002 -k gevent -w 4 -n esb --access-logfile - --error-logfile -
directory=/opt/opsany-paas/paas-ce/paas/esb
environment=PYTHONPATH=/opt/opsany-paas/.runtime/esb,BK_ENV="production",PAAS_LOGGING_DIR="/opt/opsany/logs"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/paas_esb.log
