# paas-ce/paas/examples/supervisord.d/login.ini
[program:login]
command=/opt/dev-paas/paas-runtime/login/bin/gunicorn wsgi --bind 0.0.0.0:8003 -k gevent -w 4 -n login --access-logfile - --error-logfile -
directory=/opt/dev-paas/paas-ce/paas/login
environment=DJANGO_SETTINGS_MODULE=settings,PYTHONPATH=/opt/dev-paas/paas-runtime/login,BK_ENV="production",PAAS_LOGGING_DIR="/opt/dev-paas/paas-runtime/logs"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/dev-paas/paas-runtime/logs/paas_login.log
