# paas-ce/paas/examples/supervisord.d/paas.ini
[program:paas]
command=/opt/opsany/.runtime/paas/bin/gunicorn wsgi --bind 0.0.0.0:8001 -k gevent -w 4 -n paas --access-logfile - --error-logfile -
directory=/opt/opsany-paas/paas-ce/paas/paas
environment=DJANGO_SETTINGS_MODULE=settings,PYTHONPATH=/opt/opsany/.runtime/paas,BK_ENV="production",PAAS_LOGGING_DIR="/opt/opsany/logs"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/paas_paas.log
