# paas-ce/paas/examples/supervisord.d/paasagent.ini
[program:paas-agent]
directory=/opt/opsany/paas-agent/
environment=PATH="/opt/opsany/paas-agent/bin"
command=/opt/opsany/paas-agent/bin/paas_agent
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/paas-agent-supervisord.log
