# paas-ce/paas/examples/supervisord.d/websocket.ini
[program:websocket]
command=/opt/opsany/paas-agent/apps/Envs/control/bin/daphne --proxy-headers -b 0.0.0.0 -p 8004 asgi:application
directory=/opt/opsany/paas-agent/apps/projects/control/code/control
environment=BK_ENV="testing",BK_LOG_DIR="/opt/opsany/logs",APP_ID="control"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/paas_websocket.log

