# install/conf/opsany-saas/cmdb/cmdb-supervisor.ini
[program: cmdb_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/cmdb.ini
stdout_logfile = /opt/opsany/logs/cmdb/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
stopasgroup=true
killasgroup=true

[program: cmdb_celery]
command = /usr/local/bin/python3 /opt/opsany/cmdb/manage.py celery worker -n cmdb -l INFO --autoscale=8,1
directory = /opt/opsany/cmdb
stdout_logfile = /opt/opsany/logs/cmdb/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
stopasgroup=true
killasgroup=true

[program:cmdb_beat]
command = /usr/local/bin/python3 /opt/opsany/cmdb/manage.py celery beat
directory = /opt/opsany/cmdb
stdout_logfile =/opt/opsany/logs/cmdb/celery-beat.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
stopasgroup=true
killasgroup=true