# install/conf/opsany-saas/llmops/llmops-supervisor.ini
[program: llmops_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/llmops.ini
stdout_logfile = /opt/opsany/logs/llmops/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
