# install/conf/opsany-paas/esb/esb.ini
[program:opsany-esb]
command=gunicorn wsgi --bind 0.0.0.0:8002 -k gevent -w 6 -n esb --access-logfile - --error-logfile -
directory=/opt/opsany/paas/esb
environment=BK_ENV="production",PAAS_LOGGING_DIR="/opt/opsany/logs"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/paas_esb.log
