# install/opsany-websocket/websocket.ini
[program:websocket]
command=daphne --proxy-headers -b 0.0.0.0 -p 8004 asgi:application
directory=/opt/opsany/websocket
environment=BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs",APP_ID="bastion"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/paas_websocket.log
