# install/conf/opsany-saas/repo/repo-supervisor.ini
[program: repo_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/repo.ini
stdout_logfile = /opt/opsany/logs/repo/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

