# install/conf/opsany-saas/repo/repo-supervisor.ini
[program: repo_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/repo.ini
stdout_logfile = /opt/opsany/logs/repo/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: repo_celery]
command = /usr/local/bin/python3 /opt/opsany/repo/manage.py celery worker -n repo -l INFO --autoscale=8,2
directory = /opt/opsany/repo
stdout_logfile = /opt/opsany/logs/repo/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

