# install/conf/opsany-saas/apm/apm-supervisor.ini
[program: apm_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/apm.ini
stdout_logfile = /opt/opsany/logs/apm/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
