# install/conf/opsany-saas/apm/apm-supervisor.ini
[program: apm_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/apm.ini
stdout_logfile = /opt/opsany/logs/apm/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: apm_celery]
command = /usr/local/bin/python3 /opt/opsany/apm/manage.py celery worker -n apm -l INFO --autoscale=16,4
directory = /opt/opsany/apm
stdout_logfile = /opt/opsany/logs/apm/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program:apm]
command = /usr/local/bin/python3 /opt/opsany/apm/manage.py celery beat
directory = /opt/opsany/apm
stdout_logfile =/opt/opsany/logs/apm/celery-beat.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
