# install/opsany-ce/opsany-bk-websocket/websocket.ini
[program:websocket]
command=uvicorn --proxy-headers --host 0.0.0.0 --port 8004 asgi:application --workers 4
directory=/opt/opsany/websocket
environment=BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs",APP_ID="opsany-bastion"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/paas_websocket.log
