# install/conf/opsany-saas/foot/foot-supervisor.ini
[program: foot]
command = python3 /opt/opsany/bastion-foot/maind start
stdout_logfile = /opt/opsany/logs/foot/server.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
stopasgroup=true
killasgroup=true
