# install/conf/opsany-saas/kbase/kbase-supervisor.ini
[program: kbase_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/kbase.ini
stdout_logfile = /opt/opsany/logs/kbase/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

