# install/conf/opsany-saas/monitor/monitor-supervisor.ini
[program: monitor_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/monitor.ini
stdout_logfile = /opt/opsany/logs/monitor/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

