# install/conf/opsany-saas/control/control-supervisor.ini
[program: control_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/control.ini
stdout_logfile = /opt/opsany/logs/control/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: control_celery]
command = /usr/local/bin/python3 /opt/opsany/control/manage.py celery worker -n control -l INFO --autoscale=8,2
directory = /opt/opsany/control
stdout_logfile = /opt/opsany/logs/control/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program:control_beat]
command = /usr/local/bin/python3 /opt/opsany/control/manage.py celery beat
directory = /opt/opsany/control
stdout_logfile =/opt/opsany/logs/control/celery-beat.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
