# install/conf/opsany-saas/k8s/k8s-supervisor.ini
[program: k8s_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/k8s.ini
stdout_logfile = /opt/opsany/logs/k8s/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: k8s_celery]
command = /usr/local/bin/python3 /opt/opsany/k8s/manage.py celery worker -n k8s -l INFO --autoscale=8,2
directory = /opt/opsany/k8s
stdout_logfile = /opt/opsany/logs/k8s/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program:k8s_beat]
command = /usr/local/bin/python3 /opt/opsany/k8s/manage.py celery beat
directory = /opt/opsany/k8s
stdout_logfile =/opt/opsany/logs/k8s/celery-beat.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
