# install/conf/opsany-saas/deploy/deploy-supervisor.ini
[program: deploy_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/deploy.ini
stdout_logfile = /opt/opsany/logs/deploy/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: deploy_celery]
command = /usr/local/bin/python3 /opt/opsany/deploy/manage.py celery worker -n deploy -l INFO --autoscale=8,2
directory = /opt/opsany/deploy
stdout_logfile = /opt/opsany/logs/deploy/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
