# install/conf/opsany-paas/login/login.ini
[program:opsany-login]
command=gunicorn wsgi --bind 0.0.0.0:8003 -k gevent -w 6 -n login --access-logfile - --error-logfile -
directory=/opt/opsany/paas/login
environment=DJANGO_SETTINGS_MODULE=settings,BK_ENV="production",PAAS_LOGGING_DIR="/opt/opsany/logs"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/paas_login.log
