# install/conf/opsany-saas/dashboard/dashboard-supervisor.ini
[program: dashboard_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/dashboard.ini
stdout_logfile = /opt/opsany/logs/dashboard/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
