# install/opsany-ce/opsany-stackstorm/nginx.ini
[program:paas-agent-nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/st2-nginx.log
