# install/opsany-ce/opsany-stackstorm/st2timersengine.ini
[program:st2timersengine]
command=/opt/stackstorm/st2/bin/st2timersengine --config-file /etc/st2/st2.conf
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/st2/st2-timersengine.log
