# install/conf/opsany-saas/cmp/cmp-supervisor.ini
[program: cmp_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/cmp.ini
stdout_logfile = /opt/opsany/logs/cmp/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: cmp_celery]
command = /usr/local/bin/python3 /opt/opsany/cmp/manage.py celery worker -n cmp -l INFO --autoscale=8,1
directory = /opt/opsany/cmp
stdout_logfile = /opt/opsany/logs/cmp/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program:cmp_beat]
command = /usr/local/bin/python3 /opt/opsany/cmp/manage.py celery beat
directory = /opt/opsany/cmp
stdout_logfile =/opt/opsany/logs/cmp/celery-beat.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
