# install/opsany-stackstorm/st2api.ini
[program:st2api]
command=/opt/stackstorm/st2/bin/gunicorn st2api.wsgi:application -k eventlet -b 127.0.0.1:9101 --workers 1 --threads 1 --graceful-timeout 10 --timeout 30 --log-config /etc/st2/logging.api.gunicorn.conf --error-logfile /opt/opsany/logs/st2/st2api.log
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/st2/st2-api.log
