# install/opsany-stackstorm/st2auth.ini
[program:st2auth]
command=/opt/stackstorm/st2/bin/gunicorn st2auth.wsgi:application -k eventlet -b 127.0.0.1:9100 --workers 1 --threads 1 --graceful-timeout 10 --timeout 30 --log-config /etc/st2/logging.auth.gunicorn.conf --error-logfile /opt/opsany/logs/st2/st2auth.log
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/st2/st2-auth.log
