# install/opsany-stackstorm/st2stream.ini
[program:st2stream]
command=/opt/stackstorm/st2/bin/gunicorn st2stream.wsgi:application -k eventlet -b 127.0.0.1:9102 --workers 1 --threads 10 --graceful-timeout 10 --timeout 30 --log-config /etc/st2/logging.stream.gunicorn.conf --error-logfile /var/log/st2/st2stream.log
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/st2/st2-stream.log
