# install/conf/proxy/proxy.ini
[program:opsany-proxy]
command=gunicorn wsgi --bind 0.0.0.0:8010 -k gevent -w 3 -n proxy --access-logfile - --error-logfile - --max-requests 1000 --max-requests-jitter 100
directory=/opt/opsany-proxy
environment=DJANGO_SETTINGS_MODULE=settings,OPS_ANY_ENV="production"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/proxy/proxy.log
stopsignal=TERM
stopasgroup=true
killasgroup=true


[program: opsany-proxy-celery]
command = celery --app proxy worker -l INFO -n proxy --autoscale=8,1 --max-tasks-per-child=1000
directory = /opt/opsany-proxy
environment=DJANGO_SETTINGS_MODULE=settings,OPS_ANY_ENV="production"
stdout_logfile = /opt/opsany/logs/proxy/proxy-celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
stopasgroup=true
killasgroup=true
