# install/conf/proxy/saltmaster.ini
[program:salt-master]
command=/usr/local/bin/salt-master
autostart=true
priority=10
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/proxy/salt-master.log
