# install/conf/proxy/saltapi.ini
[program:salt-api]
command=/usr/local/bin/salt-api 
autostart=true
priority=20
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/proxy/salt-api.log
