# install/conf/opsany-saas/devops/devops-supervisor.ini
[program: devops_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/devops.ini
stdout_logfile = /opt/opsany/logs/devops/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: devops_celery]
command = /usr/local/bin/python3 /opt/opsany/devops/manage.py celery worker -n devops -l INFO --autoscale=8,1
directory = /opt/opsany/devops
stdout_logfile = /opt/opsany/logs/devops/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
