# install/conf/opsany-saas/workbench/workbench-supervisor.ini
[program: workbench_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/workbench.ini
stdout_logfile = /opt/opsany/logs/workbench/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: workbench_celery]
command = /usr/local/bin/python3 /opt/opsany/workbench/manage.py celery worker -n workbench -l INFO --autoscale=8,1
directory = /opt/opsany/workbench
stdout_logfile = /opt/opsany/logs/workbench/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program:workbench_beat]
command = /usr/local/bin/python3 /opt/opsany/workbench/manage.py celery beat
directory = /opt/opsany/workbench
stdout_logfile =/opt/opsany/logs/workbench/celery-beat.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
