# install/conf/opsany-saas/pipeline/pipeline-supervisor.ini
[program: pipeline_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/pipeline.ini
stdout_logfile = /opt/opsany/logs/pipeline/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: pipeline_celery]
command = /usr/local/bin/python3 /opt/opsany/pipeline/manage.py celery worker -n pipeline -l INFO --autoscale=16,4
directory = /opt/opsany/pipeline
stdout_logfile = /opt/opsany/logs/pipeline/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

