# install/conf/opsany-saas/prom/prom-supervisor.ini
[program: prom_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/prom.ini
stdout_logfile = /opt/opsany/logs/prom/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: prom_celery]
command = /usr/local/bin/python3 /opt/opsany/prom/manage.py celery worker -n prom -l INFO --autoscale=8,2
directory = /opt/opsany/prom
stdout_logfile = /opt/opsany/logs/prom/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program:prom_beat]
command = /usr/local/bin/python3 /opt/opsany/prom/manage.py celery beat
directory = /opt/opsany/prom
stdout_logfile =/opt/opsany/logs/prom/celery-beat.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
