# install/conf/opsany-saas/auto/auto-supervisor.ini
[program: auto_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/auto.ini
stdout_logfile = /opt/opsany/logs/auto/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: auto_celery]
command = /usr/local/bin/python3 /opt/opsany/auto/manage.py celery worker -n auto -l INFO --autoscale=8,2
directory = /opt/opsany/auto
stdout_logfile = /opt/opsany/logs/auto/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program:auto_beat]
command = /usr/local/bin/python3 /opt/opsany/auto/manage.py celery beat
directory = /opt/opsany/auto
stdout_logfile =/opt/opsany/logs/auto/celery-beat.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
