# install/conf/opsany-saas/bastion/bastion-supervisor.ini
[program: bastion_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/bastion.ini
stdout_logfile = /opt/opsany/logs/bastion/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
