# install/opsany-ce/opsany-proxy/proxy.ini
[program:opsany-proxy]
command=gunicorn wsgi --bind 0.0.0.0:8010 -k gevent -w 4 -n proxy --access-logfile - --error-logfile -
directory=/opt/opsany-proxy
environment=DJANGO_SETTINGS_MODULE=settings,OPS_ANY_ENV="production"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/proxy/proxy.log

[program: opsany-proxy-celery]
command = /usr/local/bin/python3 celery --app proxy worker -l INFO -n proxy
directory = /opt/opsany-proxy
environment=DJANGO_SETTINGS_MODULE=settings,OPS_ANY_ENV="production"
stdout_logfile = /opt/opsany/logs/proxy/proxy_celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = OPS_ANY_ENV="production"
