# install/opsany-ce/opsany-appengine/appengine.ini
[program:opsany-appengine]
command=gunicorn wsgi --bind 0.0.0.0:8000 -k gevent -w 4 -n appengine --access-logfile - --error-logfile -
directory=/opt/opsany/paas/appengine
environment=BK_ENV="production",PAAS_LOGGING_DIR="/opt/opsany/logs"
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/opsany/logs/paas_appengine.log
