# install/conf/opsany-saas/job/job-supervisor.ini
[program: job_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/job.ini
stdout_logfile = /opt/opsany/logs/job/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: job_celery]
command = /usr/local/bin/python3 /opt/opsany/job/manage.py celery worker -n job -l INFO --autoscale=8,1
directory = /opt/opsany/job
stdout_logfile = /opt/opsany/logs/job/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program:job_beat]
command = /usr/local/bin/python3 /opt/opsany/job/manage.py celery beat
directory = /opt/opsany/job
stdout_logfile =/opt/opsany/logs/job/celery-beat.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
