# install/conf/opsany-saas/code/code-supervisor.ini
[program: code_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/code.ini
stdout_logfile = /opt/opsany/logs/code/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

