# install/conf/opsany-saas/event/event-supervisor.ini
[program: event_uwsgi]
command = uwsgi --ini /opt/opsany/uwsgi/event.ini
stdout_logfile = /opt/opsany/logs/event/uwsgi.log
redirect_stderr = true
autorestart = true
stopsignal = QUIT
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program: event_celery]
command = /usr/local/bin/python3 /opt/opsany/event/manage.py celery worker -n event -l INFO --autoscale=8,1
directory = /opt/opsany/event
stdout_logfile = /opt/opsany/logs/event/celery.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"

[program:event_beat]
command = /usr/local/bin/python3 /opt/opsany/event/manage.py celery beat
directory = /opt/opsany/event
stdout_logfile =/opt/opsany/logs/event/celery-beat.log
redirect_stderr = true
stopwaitsecs = 10
stopsignal = QUIT
autorestart = true
environment = BK_ENV="production",BK_LOG_DIR="/opt/opsany/logs/"
